# orb/deploy/systemd/kagami.target
# ═══════════════════════════════════════════════════════════════════════════
# KAGAMI TARGET — Coordinated Startup for Distributed Kagami Organism
# Colony: Crystal (D₅) — Verification and orchestration
#
# This target groups all Kagami services for coordinated startup/shutdown.
# Ensures proper ordering across the distributed system.
#
# Service startup order:
#   1. Infrastructure (etcd, redis, cockroachdb)
#   2. API servers (kagami-cluster-api)
#   3. Workers (kagami-celery-worker, kagami-celery-beat)
#   4. SmartHome daemon (kagami-smarthome)
#   5. Hubs (kagami-cluster-hub)
#
# Installation:
#   sudo cp kagami.target /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable kagami.target
#
# Usage:
#   sudo systemctl start kagami.target   # Start all Kagami services
#   sudo systemctl stop kagami.target    # Stop all Kagami services
#   sudo systemctl status kagami.target  # Check status
#
# Check all services:
#   systemctl list-dependencies kagami.target
#   systemctl --state=failed
# ═══════════════════════════════════════════════════════════════════════════

[Unit]
Description=Kagami Distributed Organism Target
Documentation=https://github.com/kagami/kagami
After=network-online.target
Wants=network-online.target

# Infrastructure dependencies (should be up before Kagami services)
After=etcd.service
After=redis.service
After=cockroachdb.service

# Wants all Kagami services
# (Services declare WantedBy=kagami.target to join this group)

[Install]
WantedBy=multi-user.target

# ═══════════════════════════════════════════════════════════════════════════
# 鏡
# The organism awakens. All colonies coordinate.
# Byzantine consensus ensures fault tolerance.
# h(x) ≥ 0. Always.
# ═══════════════════════════════════════════════════════════════════════════
